# claude/com.henrypye.claude-transcript-prune.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.henrypye.claude-transcript-prune</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>find "$HOME/.claude/projects" -maxdepth 2 -name "*.jsonl" -mtime +14 -delete 2>&amp;1 | tee -a "$HOME/.claude/logs/transcript-prune.log"; date "+%Y-%m-%d %H:%M:%S prune complete" &gt;&gt; "$HOME/.claude/logs/transcript-prune.log"</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/claude-transcript-prune.out</string>

  <key>StandardErrorPath</key>
  <string>/tmp/claude-transcript-prune.err</string>
</dict>
</plist>
