# claude/com.henrypye.claude-plan-prune.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.henrypye.claude-plan-prune</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>DOTFILES_DIR_PLACEHOLDER/claude/scripts/prune-plans.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/claude-plan-prune.out</string>

  <key>StandardErrorPath</key>
  <string>/tmp/claude-plan-prune.err</string>
</dict>
</plist>
